# scripts/bcl-umj-probe.service
[Unit]
Description=BCL UMJ Probe FastAPI Service
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
WorkingDirectory=/usr/local/bcl-umj-probe
Environment="PATH=/usr/local/bcl-umj-probe/venv/bin:/usr/bin:/bin"
ExecStart=/usr/local/bcl-umj-probe/venv/bin/uvicorn app:api --host 0.0.0.0 --port 8000 --workers 4 --log-level info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
